# DiskUsage/cron/diskusage.cron
#
# Every night run disk_usage.pl and write a cache for the diskusage
# reporting interface.
#
1 1 * * * root /usr/bin/disk_usage.pl -i /var/www/domains/gsc.wustl.edu/diskusage/cgi-bin/du.cache -D /gsc/scripts/share/automount/gscmnt >/dev/null 2>&1
